# homework6_Systemd/1/lib/systemd/system/grep_log.timer
[Unit]
Description=timer for grep_log.service
 
[Timer]
OnCalendar=*:*:0,30
Unit=grep_log.service

[Install]
WantedBy=timers.target
